# systemd/ray-worker.service
[Unit]
Description=Ray worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/home/%i/dgx-cluster-check/.env
ExecStart=/home/%i/dgx-cluster-check/venv/bin/ray start --address=${MASTER_ADDR}:${MASTER_PORT}
ExecStop=/home/%i/dgx-cluster-check/venv/bin/ray stop
Restart=on-failure

[Install]
WantedBy=multi-user.target

