# systemd/ray-head.service
[Unit]
Description=Ray head node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/home/%i/dgx-cluster-check/.env
ExecStart=/home/%i/dgx-cluster-check/venv/bin/ray start --head --node-ip-address=%h --port=${MASTER_PORT}
ExecStop=/home/%i/dgx-cluster-check/venv/bin/ray stop
Restart=on-failure

[Install]
WantedBy=multi-user.target

